# keepalive.service
[Unit]
Description=Adaptive Resource Workload Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /root/keepalive.py
Restart=always
RestartSec=15
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
CPUQuota=360%
MemoryMax=35%
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectSystem=full
ReadWritePaths=/var/lib/keepalive /etc/keepalive /var/log/keepalive
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
